# systemd/eda-modbus-bridge.service
[Unit]
Description=HTTP bridge for Enervent ventilation units with EDA automation

[Service]
Type=simple
EnvironmentFile=/etc/default/eda-modbus-bridge
ExecStart=/usr/bin/node /opt/eda-modbus-bridge/dist/eda-modbus-bridge.js $OPTIONS
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
